# root.crontab
* * * * * /etc/init.d/20fw >> /data/tmp/crontest